# deploy/iot-web.service
[Unit]
Description=IoT Web Dashboard
After=network.target

[Service]
User=gruppe1
WorkingDirectory=/home/gruppe1/IoT/web
ExecStart=/usr/bin/python3 -m http.server 3000 --bind 0.0.0.0
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target